# container/etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=service nginx restart
